# docker/etc/supervisor/conf.d/srf-component-library.conf
[program:srf-component-library]
directory=/app
command=/usr/bin/npm start
autostart=true
autorestart=true
startretries=3
startsecs=2
redirect_stdout=true
redirect_stderr=true